# aurora/units/aurora-health.timer
[Unit]
Description=Aurora 2.1 - Timer do health monitor (a cada 30min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min

[Install]
WantedBy=timers.target
