# aurora/units/browser.slice
[Unit]
Description=Aurora 2.1 - Slice para browsers (Chrome, Firefox, Brave)
Before=slices.target

[Slice]
MemoryHigh=8G
MemoryMax=10G
MemorySwapMax=4G
CPUWeight=80
