# aurora/units/claude.slice
[Unit]
Description=Aurora 2.2 - Slice para Claude Code (limites de memoria + prioridade alta)

[Slice]
MemoryHigh=10G
MemoryMax=12G
CPUWeight=200
IOWeight=200
