# aurora/units/electron.slice
[Unit]
Description=Aurora 2.1 - Slice para apps Electron (Slack, Discord, VS Code, Zoom)
Before=slices.target

[Slice]
MemoryHigh=4G
MemoryMax=6G
MemorySwapMax=2G
CPUWeight=70
