# aurora/units/ollama.slice
[Unit]
Description=Aurora 2.1 - Slice para Ollama (limites de memoria RAM)
Before=slices.target

[Slice]
MemoryHigh=6G
MemoryMax=9G
MemorySwapMax=4G
CPUWeight=80
